# zhcp/bin/IUCV/iucvserd.service
[Unit]
Description=IUCV agent service which is used to communicate with IUCV client in OPNCLOUD.
Wants=
After=cloud-init-local.service
Before=

[Service]
Type=simple
ExecStart=/usr/bin/iucvserv start
StandardOutput=journal+console

[Install]
WantedBy=multi-user.target
